# task5/mysql-backup.service
[Unit]
Description=MySQL Backup Service

[Service]
Type=oneshot
ExecStart=/path/to/your/backup-script.sh
